# ansible-slurm.path
[Unit] 
Description=Ansible slurm.conf file monitor

[Path]
PathModified=/ansible/staging/slurm.conf

[Install]
WantedBy=multi-user.target

