# ansible-hosts-allow.path
[Unit] 
Description=Ansible hosts.allow file monitor

[Path]
PathModified=/ansible/staging/hosts.allow

[Install]
WantedBy=multi-user.target

