# ansible-hosts.path
[Unit] 
Description=Ansible hosts file monitor

[Path]
PathModified=/ansible/staging/hosts

[Install]
WantedBy=multi-user.target

